# docker/cron/cleanup.cron
*/5 * * * * /usr/local/bin/cleanup-old-files.sh >> /var/log/cron.log 2>&1